# root/etc/s6-overlay/s6-rc.d/init-mod-universal-borgmatic-postinstall/run
#!/usr/bin/with-contenv bash
echo "Validating configuration..."

# permissions
lsiown -R abc:abc \
    /app \
    /config

chmod 600 /config/borgmatic.yaml
chmod 600 /config/crontabs/abc

borgmatic config validate --config /config/borgmatic.yaml

export BORG_SOURCE_DIRECTORY="${BORG_SOURCE_DIRECTORY:=/config}"
export BORG_REPO="${BORG_REPO:=/tmp/borgmatic-backup}"

mkdir -p BORG_SOURCE_DIRECTORY
mkdir -p BORG_REPO

if [ -z ${BORG_REPO_KEY+x} ]; 
then 
  echo "No repo key set; repo will be unencrypted; 
  borgmatic init --config /config/borgmatic.yaml
else 
  echo "Repo key set; repo will be encrypted"; 
  borgmatic init --encryption repokey-aes-ocb --config /config/borgmatic.yaml
fi
