# root/etc/s6-overlay/s6-rc.d/init-mod-universal-borgmatic-postinstall/run
#!/usr/bin/with-contenv bash

mkdir -p ${BORGMATIC_ROOT}/.config/borgmatic/ ${BORGMATIC_ROOT}/.config/borg/

# copy default config files if they don't exist
if [[ ! -f ${BORGMATIC_ROOT}/.config/borgmatic/config.yaml ]]; then
    cp /defaults/config.yaml.sample ${BORGMATIC_ROOT}/.config/borgmatic/config.yaml
fi

# rather than deal with root crontabs and their possible conflicts with other containers (i.e, swag),
# we use alpine's /etc/periodic run-parts scripts
# periodic can be: 15min, hourly, daily, weekly, or monthly
cp /app/borgmatic-backup /etc/periodic/${BORGMATIC_PERIODIC}
chmod +rx /etc/periodic/${BORGMATIC_PERIODIC}/borgmatic-backup

# permissions
lsiown -R abc:abc /config

chmod +r ${BORGMATIC_ROOT}/.config/borgmatic/config.yaml




echo "Validating configuration..."

borgmatic config validate --config ${BORGMATIC_ROOT}/.config/borgmatic/config.yaml

echo "Borg repo encryption mode set to ${BORGMATIC_ENCRYPTION}"; 
borgmatic rcreate --encryption ${BORGMATIC_ENCRYPTION} --config ${BORGMATIC_ROOT}/.config/borgmatic/config.yaml --make-parent-dirs
