# root/etc/s6-overlay/s6-rc.d/init-mod-universal-borgmatic-postinstall/run
#!/usr/bin/with-contenv bash

export BORG_SOURCE_DIRECTORY="${BORG_SOURCE_DIRECTORY:=/config}"
export BORG_REPO="${BORG_REPO:=/tmp/borgmatic-backup}"

mkdir -p ${BORG_SOURCE_DIRECTORY}
mkdir -p ${BORG_REPO}

chmod 750 -R ${BORG_REPO}

lsiown abc:abc -R ${BORG_REPO}

echo "Validating configuration..."

borgmatic config validate --config /config/borgmatic.yaml

if [ -z ${BORG_REPO_KEY+x} ]; 
then 
  echo "No repo key set; repo will be unencrypted"; 
  borgmatic init --encryption none --config /config/borgmatic.yaml
else 
  echo "Repo key set; repo will be encrypted"; 
  borgmatic init --encryption repokey-aes-ocb --config /config/borgmatic.yaml
fi
