# root/etc/s6-overlay/s6-rc.d/init-mod-universal-borgmatic-postinstall/run
#!/usr/bin/with-contenv bash

export BORG_ENCRYPTION="${BORG_ENCRYPTION:=none}"

echo "Validating configuration..."

borgmatic config validate --config /config/borgmatic.yaml

echo "Borg repo encryption mode set to ${BORG_ENCRYPTION}"; 
borgmatic rcreate --encryption ${BORG_ENCRYPTION} --config /config/borgmatic.yaml --make-parent-dirs
