# root/defaults/crontabs/root
# You can drop this file into /etc/cron.d/ to run borgmatic nightly.

*/5 * * * * borgmatic --verbosity -1 --syslog-verbosity 1 --config /config/borgmatic.yaml
