# root/etc/s6-overlay/s6-rc.d/init-mod-universal-borgmatic/run
#!/usr/bin/with-contenv bash

# default environment variables
echo "${BORGMATIC_ENCRYPTION:-none}" > /run/s6/container_environment/BORGMATIC_ENCRYPTION
echo "${BORGMATIC_PERIODIC:-daily}" > /run/s6/container_environment/BORGMATIC_PERIODIC
echo "${BORGMATIC_ROOT:-/config}" > /run/s6/container_environment/BORGMATIC_ROOT
echo "${BORGMATIC_SOURCE:-/config}" > /run/s6/container_environment/BORGMATIC_SOURCE
echo "${BORGMATIC_REPO:-/tmp/borgmatic-backup}" > /run/s6/container_environment/BORGMATIC_REPO

# "${BORGMATIC_RESTORE:=latest}"

# install packages
if [ -f /usr/bin/apt ]; then
    ## Ubuntu
    echo "\
        borgbackup \
        borgmatic \
        openssh-client" >> /mod-repo-packages-to-install.list
fi

if [ -f /sbin/apk ]; then
    # Alpine
    echo "\
        borgbackup \
        borgmatic \
        openssh" >> /mod-repo-packages-to-install.list
fi
